# 1.0/tt-rss.sv.conf
[program:tt-rss]
command=/usr/share/tt-rss/www/update.php --daemon
numprocs=1
autostart=true
autorestart=true
user=www-data
